# docker/supervisor/production-no-docker.conf
; ==========================================================
; Supervisor config for production WITHOUT Docker
; Copy this file to /etc/supervisor/conf.d/icc-munich.conf
; Then run: sudo supervisorctl reread && sudo supervisorctl update
; ==========================================================

[program:icc-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/icc-munich/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/icc-munich
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/icc-munich/storage/logs/queue-worker.log
stopwaitsecs=3600

[program:icc-scheduler]
process_name=%(program_name)s
command=sh -c "while true; do php /var/www/icc-munich/artisan schedule:run --verbose --no-interaction; sleep 60; done"
directory=/var/www/icc-munich
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/icc-munich/storage/logs/scheduler.log
stopwaitsecs=60
